# systemd/user/change_wallpaper.service
[Unit]
Description=Change wallpaper using SWWW
After=swww.service

[Service]
ExecStart=/bin/bash change-wallpaper

[Install]
WantedBy=default.target
